# chatty.conf
[program:chatty]
command=/.cabal/bin/chatty -p 80
directory=/opt/chatty
autostart=True
autorestart=True
redirect_stderr=True